# UI/systemd/raspi-ui.service
[Unit]
Description=Photoelectric Current Experiment UI (1024x600 fullscreen)
After=display-manager.service
Wants=display-manager.service

[Service]
Type=simple
User=cc
Group=cc
WorkingDirectory=/home/cc/Desktop/UICode
Environment=HOME=/home/cc
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash /home/cc/Desktop/UICode/UI/run_hardware_ui.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=graphical.target
